# meta-custom-image-homeassistant/recipes-monitoring/docker-system-monitor/files/docker-system-monitor.service
[Unit]
Description=Docker System Monitor Service
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
ExecStart=/usr/bin/docker-system-monitor
Restart=always
RestartSec=30
User=root
StandardOutput=journal
StandardError=journal

# Environment
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log /tmp /var/run/docker.sock

# Docker socket access
SupplementaryGroups=docker

[Install]
WantedBy=multi-user.target 